# argonaut-agents/debian/argonaut-agents-printmanager.cron.d
* * * * *	root /usr/sbin/argonaut-enable-printer
